# deploy/systemd/tradebot-dashboard.service
# systemd USER unit — Linux/Jetson equivalent of the macOS launchd
# dashboard agent (deploy/launchd/com.tradebot.dashboard.plist).
#
# Installed per-user (no sudo) via:
#   bash scripts/tradebotctl.sh dashboard-install
# which copies this template to ~/.config/systemd/user/, rewrites the
# __TRADEBOT_ROOT__ / __TRADEBOT_PY__ placeholders, then:
#   systemctl --user daemon-reload
#   systemctl --user enable --now tradebot-dashboard
#
# Default bind: 127.0.0.1:8000 (see scripts/run_dashboard.py).
[Unit]
Description=tradebot dashboard (FastAPI read-only UI on localhost:8000)
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
WorkingDirectory=__TRADEBOT_ROOT__
Environment=TZ=America/New_York
ExecStart=__TRADEBOT_PY__ __TRADEBOT_ROOT__/scripts/run_dashboard.py
KillSignal=SIGTERM
TimeoutStopSec=10
Restart=on-failure
RestartSec=5
StandardOutput=append:__TRADEBOT_ROOT__/logs/dashboard.out
StandardError=append:__TRADEBOT_ROOT__/logs/dashboard.err

[Install]
WantedBy=default.target
