# deploy/systemd/tradebot-watchdog.service
# systemd USER unit — Linux/Jetson equivalent of the macOS launchd
# watchdog (deploy/launchd/com.tradebot.paper.plist).
#
# This runs the watchdog wrapper (scripts/watchdog_run.py), which:
#   - spawns run_paper.py as a child
#   - fires a Discord/Slack alert on crash (via src.notify)
#   - records every start/exit to logs/watchdog_events.jsonl
#   - restarts the child if heartbeat (logs/heartbeat.txt) goes stale
#
# Installed per-user (no sudo) via:
#   bash scripts/tradebotctl.sh watchdog-install
# which copies this template to ~/.config/systemd/user/, rewrites the
# __TRADEBOT_ROOT__ / __TRADEBOT_PY__ placeholders, then:
#   systemctl --user daemon-reload
#   systemctl --user enable --now tradebot-watchdog
#
# To survive logout (stay running after SSH session closes):
#   sudo loginctl enable-linger $USER
[Unit]
Description=tradebot watchdog (supervises run_paper.py)
After=network-online.target
Wants=network-online.target
# Restart-storm protection: if the *watchdog itself* crashes 5 times in
# 60s, stop trying. The child-process crash loop is handled inside
# watchdog_run.py with its own backoff.
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
WorkingDirectory=__TRADEBOT_ROOT__
# .env is loaded via python-dotenv inside the Python entrypoints, so we
# don't need EnvironmentFile= here. TZ ensures the bot's session clock
# uses US Eastern regardless of host locale.
Environment=TZ=America/New_York
Environment=TRADEBOT_PY=__TRADEBOT_PY__
ExecStart=__TRADEBOT_PY__ __TRADEBOT_ROOT__/scripts/watchdog_run.py
# Watchdog handles clean shutdown of its child; 30s is plenty.
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=on-failure
RestartSec=10
StandardOutput=append:__TRADEBOT_ROOT__/logs/watchdog.out
StandardError=append:__TRADEBOT_ROOT__/logs/watchdog.err

[Install]
WantedBy=default.target
