# deploy/launchd/com.tradebot.scan_premarket.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Pre-market ticker scan. Fires at 08:30 ET weekdays (1 hour before
  open). Posts to Discord title='scan_premarket'. Uses Groq 70B
  (cloud) for commentary when GROQ_API_KEY is set, else local Ollama.

  Install:  scripts/tradebotctl.sh scan-premarket-install
  Status:   launchctl list | grep tradebot.scan_premarket
  Logs:     logs/scan_premarket.out / .err
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tradebot.scan_premarket</string>

  <key>WorkingDirectory</key>
  <string>__TRADEBOT_ROOT__</string>

  <key>ProgramArguments</key>
  <array>
    <string>__TRADEBOT_PY__</string>
    <string>__TRADEBOT_ROOT__/scripts/run_morning_scan.py</string>
    <string>--session</string>
    <string>pre</string>
    <string>--limit</string>
    <string>12</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    <key>TZ</key>
    <string>America/New_York</string>
  </dict>

  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <key>StandardOutPath</key>
  <string>__TRADEBOT_ROOT__/logs/scan_premarket.out</string>
  <key>StandardErrorPath</key>
  <string>__TRADEBOT_ROOT__/logs/scan_premarket.err</string>
</dict>
</plist>
