# deploy/launchd/com.tradebot.news_collector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Continuous market-news collector. Pulls headlines from Finnhub,
  Yahoo, Alpaca, SEC EDGAR, Fed/Treasury RSS, plus optional
  MarketAux/AlphaVantage/NewsAPI. Posts branded embeds to Discord
  channel set via DISCORD_WEBHOOK_URL_NEWS.

  Install:  scripts/tradebotctl.sh news-install
  Status:   launchctl list | grep tradebot.news_collector
  Logs:     logs/news.out / .err
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tradebot.news_collector</string>

  <key>WorkingDirectory</key>
  <string>__TRADEBOT_ROOT__</string>

  <key>ProgramArguments</key>
  <array>
    <string>__TRADEBOT_PY__</string>
    <string>__TRADEBOT_ROOT__/scripts/news_collector.py</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    <key>TZ</key>
    <string>America/New_York</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>StandardOutPath</key>
  <string>__TRADEBOT_ROOT__/logs/news.out</string>
  <key>StandardErrorPath</key>
  <string>__TRADEBOT_ROOT__/logs/news.err</string>
</dict>
</plist>
