# deploy/jetson/services/tradebot-dashboard.service
# Separate systemd unit for the read-only dashboard.
# Bound to 127.0.0.1 by default; access via SSH tunnel:
#   ssh -L 8000:localhost:8000 orin@<jetson-ip>
[Unit]
Description=tradebot dashboard — Jetson
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=orin
Group=orin
WorkingDirectory=/home/orin/tradebot
EnvironmentFile=/home/orin/tradebot/.env
Environment=TZ=America/New_York
ExecStart=/home/orin/tradebot/.venv/bin/python scripts/run_dashboard.py --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/orin/tradebot/logs/dashboard.out
StandardError=append:/home/orin/tradebot/logs/dashboard.err

# Let this share the 4 efficiency cores so it can't starve the trade loop.
CPUAffinity=8 9 10 11

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/home/orin/tradebot/logs

[Install]
WantedBy=multi-user.target
